# crontab
0 */12 * * * python3 INSERT FILE PATH HERE
